# recipes-core/vacsystemupdater/vacsystemupdater-1.0/vacsystemupdater.service
[Unit]
Description=Starts system updater service
Before=NetworkManager.service
Requires=fb1.service
After=fb1.service systemd-udev-trigger.service systemd-udevd.service
# Requires=psplash-vacmonitor.service
After=psplash-vacmonitor.service
DefaultDependencies=no

[Service]
User=root
Type=oneshot
KillMode=none
Environment=USERHOMEDIR=/home/@ORIGINAL_USER@
ExecStart=/usr/sbin/vacsystemupdater.sh

[Install]
WantedBy=sysinit.target


